# systemd/certbot-s3.timer
[Unit]
Description=certbot-s3

[Timer]
OnCalendar=monthly
AccuracySec=1min
Persistent=true

[Install]
WantedBy=multi-user.target
